# deploy/systemd/alatube.service
[Unit]
Description=AlaTube Go API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alatube
Group=alatube
WorkingDirectory=/opt/alatube
EnvironmentFile=/etc/alatube/alatube.env
ExecStart=/opt/alatube/bin/alatube
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/alatube /var/cache/alatube
StateDirectory=alatube
CacheDirectory=alatube
RuntimeDirectory=alatube

CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateDevices=true
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

CPUQuota=100%
MemoryMax=768M
TasksMax=128
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target

